# macos/BubuQuotaPanel/Resources/io.github.mayday-materials.bubu-quota-panel.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.mayday-materials.bubu-quota-panel</string>
  <key>ProgramArguments</key>
  <array>
    <string>__EXECUTABLE__</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>BUBU_PANEL_HEALTH_FILE</key>
    <string>__HEALTH_PATH__</string>
    <key>BUBU_SHOW_MARKET_PRICES</key>
    <string>true</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
  <key>ThrottleInterval</key>
  <integer>2</integer>
  <key>StandardOutPath</key>
  <string>__LOG_PATH__</string>
  <key>StandardErrorPath</key>
  <string>__LOG_PATH__</string>
</dict>
</plist>
